# deployment/supervisor/nurisk-queue.conf
; /etc/supervisor/conf.d/nurisk-queue.conf
; sudo supervisorctl reread && sudo supervisorctl update
;
; 2 worker processes:
;   - pdf-generation: prioritas tinggi, 2 worker
;   - default: prioritas rendah, 1 worker
;
; Graceful shutdown: worker diberi waktu 300s untuk selesaikan job aktif.

[program:nurisk-queue-pdf]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/nurisk/artisan queue:work database
    --queue=pdf-generation
    --tries=3
    --backoff=30
    --sleep=3
    --timeout=120
    --memory=256
numprocs=2
autostart=true
autorestart=true
startretries=5
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/nurisk/storage/logs/queue-pdf.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=300

[program:nurisk-queue-default]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/nurisk/artisan queue:work database
    --queue=default
    --tries=3
    --backoff=10
    --sleep=3
    --timeout=60
    --memory=256
numprocs=1
autostart=true
autorestart=true
startretries=5
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/nurisk/storage/logs/queue-default.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=60
